# pkg.debian/compfacts-service.cron.d
0 6 * * * www-compfacts /usr/sbin/compfacts-backup
