# pkg.debian/compfacts-service.cron.d
0 */4 * * * www-compfacts /usr/sbin/compfacts-backup
